# .docker/php/crontab
* * * * *  /usr/local/bin/create-tenant-db.sh >> /var/log/cron.log 2>&1
